# pi/systemd/moki-touch-fix.service
[Unit]
Description=Fix Goodix touchscreen probe (WM8960 I2C timing workaround)
DefaultDependencies=no
After=sysinit.target
Before=moki-native.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'DEV=/sys/bus/i2c/drivers/Goodix-TS/10-005d; BIND=/sys/bus/i2c/drivers/Goodix-TS/bind; if [ -e "$DEV" ]; then echo "Goodix touchscreen already bound"; exit 0; fi; if [ ! -w "$BIND" ]; then echo "Goodix bind path unavailable: $BIND"; exit 0; fi; for i in 1 2 3 4 5; do echo "Goodix bind attempt $i"; if echo 10-005d > "$BIND" 2>/dev/null; then sleep 1; if [ -e "$DEV" ]; then echo "Goodix touchscreen bound"; exit 0; fi; fi; sleep 1; done; echo "Goodix touchscreen bind failed after retries"; exit 0'

[Install]
WantedBy=multi-user.target
